# workshop/supervisor/octant.conf
[program:octant]
process_name=octant
command=/home/eduk8s/bin/octant
env=OCTANT_LISTENER_ADDR=0.0.0.0:7777
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
